# glowkey.service
[Unit]
Description=GlowKey Backlight Restore
Documentation=https://github.com/joaomjbraga/glowkey
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=@GLOWKEY_PATH@ restore
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
